# templates/cntlm.service.j2
{{ ansible_managed | comment }}
[Unit]
Description=cntlm
After=network.target

[Service]
Type=simple

ExecStart=/usr/sbin/cntlm -c /etc/cntlm.conf -f

[Install]
WantedBy=multi-user.target
